# Services/Transmission/transmission.conf
[program:samba]
command = bash /home/pi/Documents/pi-server/Services/Transmission/runner.sh
autostart = true
autorestart = true