# Services/Plex/plex.conf
[program:plex]
command = bash /home/pi/Documents/pi-server/Services/Plex/runner.sh
autostart = true
autorestart = true