# Docker/Plex/plex.conf
[program:plex]
command = /home/pi/Documents/pi-server/Docker/Plex/runner.sh
autostart = true
autorestart = true