# Docker/Transmission/transmission.conf
[program:samba]
command = /home/pi/Documents/pi-server/Docker/Transmission/runner.sh
autostart = true
autorestart = true