# Docker/OwnCloud/owncloud.conf
[program:owncloud]
command = /home/pi/Docker/OwnCloud/runner.sh
autostart = true
autorestart = true