# Services/OwnCloud/owncloud.conf
[program:owncloud]
command = bash /home/pi/Documents/pi-server/Docker/OwnCloud/runner.sh
autostart = true
autorestart = true