# Docker/Samba/samba.conf
[program:samba]
command = /home/pi/Documents/pi-server/Docker/Samba/runner.sh
autostart = true
autorestart = true